# daemon/ce-image-server.service
[Unit]
Description=Daemon for image server

[Service]
Type=simple
ExecStart=/usr/bin/ce.image.server.daemon.o

[Install]
WantedBy=multi-user.target